# vp-edge/rootfs/etc/services.d/vp-edge/finish
#!/usr/bin/with-contenv bashio
bashio::log.warning "VoltPure Edge gestopt — exit-code ${1}"
